# roles/native-docker-compose-restart-unhealthy/files/docker-compose-restart-unhealthy.timer
[Unit]
Description=starts docker-compose-restart-unhealthy.service

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
