# roles/native-btrfs-health-check/templates/btrfs-health-check.service.j2
[Unit]
Description=Check btrfs status
OnFailure=systemd-email@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash {{docker_btrfs_health_check_folder}}btrfs-health-check.sh
