# roles/monitor-bot-btrfs/templates/monitor-bot-btrfs.service.j2
[Unit]
Description=Check btrfs status
OnFailure=alert-core.cymais@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash {{docker_health_btrfs_folder}}monitor-bot-btrfs.sh
