# roles/heal-docker/templates/heal-docker.service.j2
[Unit]
Description=restart unhealthy docker containers
OnFailure=systemd-notifier@%n.service

[Service]
Type=oneshot
ExecStart=/bin/python {{heal_docker}}heal-docker.py
