# roles/heal-docker/templates/heal-docker.timer.j2
[Unit]
Description=starts heal-docker.service

[Timer]
OnCalendar={{on_calendar_docker_compose_restart_unhealthy}}
RandomizedDelaySec={{randomized_delay_sec}}
Persistent=false

[Install]
WantedBy=timers.target
