# roles/sys-hlth-journalctl/templates/sys-hlth-journalctl.service.j2
[Unit]
Description=checking journalctl health
OnFailure=sys-alm-compose.cymais@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash {{health_journalctl_folder}}sys-hlth-journalctl.sh
