# roles/server_native-docker-health-check/templates/docker-health-check.service.j2
[Unit]
Description=Checking docker health
OnFailure=systemd-email@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash {{docker_health_check_folder}}docker-health-check.sh
