# roles/native-backups-cleanup/files/backups-cleanup.timer
[Unit]
Description=starts backups-cleanup.service

[Timer]
OnCalendar=19:00
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
