# roles/pc_application-msi-keyboard-color/templates/keyboard-color.service.j2
[Unit]
Description=Keyboard Color Service

[Service]
Type=oneshot
ExecStart=/opt/keyboard_color.sh {{ vendor_and_product_id }}