# roles/sys-rpr-btrfs-blnc/templates/sys-rpr-btrfs-blnc.service.j2
[Unit]
Description=auto balance btrfs
OnFailure=sys-alm-compose.infinito@%n.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'btrfs-auto-balancer 90 10'