# roles/health-nginx/templates/health-nginx.service.j2
[Unit]
Description=Check nginx configuration status
OnFailure=systemd-notifier@%n.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 {{ nginx_health_check_folder }}health-nginx.py
