# roles/sys-hlth-docker-container/templates/sys-hlth-docker-container.service.j2
[Unit]
Description=Checking docker health
OnFailure=sys-alm-compose.cymais@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash {{health_docker_container_folder}}sys-hlth-docker-container.sh
