# roles/backup-docker-to-local/templates/backup-docker-to-local.service.j2
[Unit]
Description=docker volume backup
OnFailure=systemd-notifier@%n.service cleanup-failed-docker-backups.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python {{docker_volume_backup_folder}}backup-docker-to-local.py