# roles/native-btrfs-health-check/files/btrfs-health-check.timer
[Unit]
Description=starts btrfs-health-check.service

[Timer]
OnCalendar=12:00

[Install]
WantedBy=timers.target
