# roles/docker-mailu/templates/deploy-letsencrypt-mailu.service.j2
[Unit]
Description=Let's Encrypt Mailu Deploy
OnFailure=systemd-notifier.cymais@%n.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash {{path_administrator_scripts}}mailu/deploy-letsencrypt-mailu.sh
