# roles/docker-mailu/templates/deploy-letsencrypt-mailu.service.j2
[Unit]
Description=Let's Encrypt Mailu Deploy
OnFailure=systemd-email@%n.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /usr/local/bin/mailu/deploy-letsencrypt-mailu.sh
