# roles/system-btrfs-auto-balancer/templates/system-btrfs-auto-balancer.service.j2
[Unit]
Description=auto balance btrfs
OnFailure=systemd-notifier.cymais@%n.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'btrfs-auto-balancer 90 10'