# roles/net-wireguard-plain/files/set-mtu.service
[Unit]
Description=set MTU 
Before=wg-quick@wg0.infinito.service

[Service]
Type=oneshot
ExecStart=bash /usr/local/bin/set-mtu.sh

[Install]
RequiredBy=wg-quick@wg0.infinito.service