# roles/update-docker/templates/update-docker.service.j2
[Unit]
Description=Updates Docker Instances
OnFailure=systemd-notifier@%n.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c '/usr/bin/python {{update_docker_script}} {{path_docker_compose_files}}'