# roles/native-pull-primary-backups/templates/pull-remote-backups.timer
[Unit]
Description=starts pull remote backup timer

[Timer]
OnCalendar=22:00
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
