# roles/native-pull-primary-backups/templates/pull-remote-backups.timer.j2
[Unit]
Description=starts pull remote backup timer

[Timer]
OnCalendar={{on_calendar_pull_primary_backups}}

[Install]
WantedBy=timers.target
