# roles/native-free-disc-space/templates/free-disc-space.service.j2
[Unit]
Description=freeing disc space
OnFailure=systemd-email@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash {{free_disc_space_folder}}free-disc-space.sh {{size_percent_free_disc_space}}
