# roles/native-journalctl-health-check/files/journalctl-health-check.timer
[Unit]
Description=starts journalctl-health-check.service

[Timer]
OnCalendar=00:00

[Install]
WantedBy=timers.target
