# roles/mon-bot-msmtp/templates/mon-bot-msmtp.service.j2
[Unit]
Description=Check msmtp liveliness
OnFailure=alert-telegram.cymais@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash {{ health_msmtp_folder }}mon-bot-msmtp.sh
