# roles/maint-btrfs-auto-balancer/templates/maint-btrfs-auto-balancer.service.j2
[Unit]
Description=auto balance btrfs
OnFailure=alert-compose.cymais@%n.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'btrfs-auto-balancer 90 10'