# roles/health-docker/templates/health-docker.timer.j2
[Unit]
Description=starts health-docker.service

[Timer]
OnCalendar={{on_calendar_health_docker}}
RandomizedDelaySec={{randomized_delay_sec}}
Persistent=false

[Install]
WantedBy=timers.target
