# roles/health-docker/templates/health-docker.service.j2
[Unit]
Description=Checking docker health
OnFailure=systemd-notifier@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash {{health_docker_folder}}health-docker.sh
