# roles/pc_application_caffeine/templates/caffeine.service.j2
[Unit]
Description=Caffeine auto start

[Service]
ExecStart={{ caffeine_executable_path }}

[Install]
WantedBy=multi-user.target
