# roles/systemd_notifier_telegram/templates/systemd_notifier_telegram@.service.j2
[Unit]
Description=status Telegram message for %i to user

[Service]
Type=oneshot
ExecStart=/bin/bash {{ systemd_telegram_script }} %i
User=root
Group=systemd-journal
