# roles/native-backups-cleanup/templates/backups-cleanup.timer.j2
[Unit]
Description=starts backups-cleanup.service

[Timer]
OnBootSec=1min
OnCalendar={{on_calendar_backups_cleanup}}

[Install]
WantedBy=timers.target
