# roles/native-docker-volume-backup/templates/docker-volume-backup.timer.j2
[Unit]
Description=backups the docker volumes

[Timer]
OnCalendar={{on_calendar_docker_volume_backup}}
RandomizedDelaySec={{randomized_delay_sec}}
Persistent=false

[Install]
WantedBy=timers.target
