# roles/native-docker-volume-backup/templates/docker-volume-backup.timer
[Unit]
Description=backups the docker volumes

[Timer]
OnCalendar=3:00
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
