# roles/monitor-bot-docker-container/templates/monitor-bot-docker-container.service.j2
[Unit]
Description=Checking docker health
OnFailure=alert-compose.cymais@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash {{health_docker_container_folder}}monitor-bot-docker-container.sh
