# roles/sys-hlth-msmtp/templates/sys-hlth-msmtp.service.j2
[Unit]
Description=Check msmtp liveliness
OnFailure=sys-alm-telegram.infinito@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash {{ health_msmtp_folder }}sys-hlth-msmtp.sh
