# roles/free-disc-space/templates/free-disc-space.timer.j2
[Unit]
Description=starts free-disc-space.service

[Timer]
OnCalendar={{on_calendar_free_disc_space}}
RandomizedDelaySec={{randomized_delay_sec}}
Persistent=true

[Install]
WantedBy=timers.target
