# roles/sys-hlth-webserver/templates/sys-hlth-webserver.service.j2
[Unit]
Description=Check nginx configuration status
OnFailure=sys-alm-compose.infinito@%n.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 {{ health_nginx_folder }}sys-hlth-webserver.py
