# roles/native-journalctl-health-check/templates/journalctl-health-check.timer.j2
[Unit]
Description=starts journalctl-health-check.service

[Timer]
OnCalendar={{on_calendar_journalctl_health_check}}

[Install]
WantedBy=timers.target
