# roles/monitor-bot-journalctl/templates/monitor-bot-journalctl.service.j2
[Unit]
Description=checking journalctl health
OnFailure=alert-core.cymais@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash {{health_journalctl_folder}}monitor-bot-journalctl.sh
