# client-playbook/roles/application-wireguard/files/set-mtu.service
[Unit]
Description=set MTU 
Before=wg-quick@wg0.service

[Service]
Type=oneshot
ExecStart=set-mtu.sh

[Install]
RequiredBy=wg-quick@wg0.service