# roles/sys-hlth-btrfs/templates/sys-hlth-btrfs.service.j2
[Unit]
Description=Check btrfs status
OnFailure=sys-alm-compose.cymais@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash {{docker_health_btrfs_folder}}sys-hlth-btrfs.sh
