# roles/server_native-docker-volume-backup/templates/docker-volume-backup.service.j2
[Unit]
Description=docker volume backup
OnFailure=systemd-email@%n.service

[Service]
Type=oneshot
ExecStart=/usr/bin/python {{docker_volume_backup_folder}}server_docker-volume-backup.py