# roles/alert-compose/templates/alert-compose@.service.j2
[Unit]
Description=Notifier for %i

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl start alert-telegram.cymais@%i.service alert-email.cymais@%i.service
User=root
Group=systemd-journal
