# roles/health-journalctl/templates/health-journalctl.service.j2
[Unit]
Description=checking journalctl health
OnFailure=systemd-notifier.cymais@%n.service

[Service]
Type=oneshot
ExecStart=/bin/bash {{health_journalctl_folder}}health-journalctl.sh
