# roles/independent_backup-to-swappable/templates/backup-to-swappable.service.j2
[Unit]
Description=backing up data
OnFailure=systemd-email@%n.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash {{backup_to_swappable_script_path}}