# roles/systemd_notifier_email/templates/systemd_notifier_email@.service.j2
[Unit]
Description=status email for %i to user

[Service]
Type=oneshot
ExecStart=/bin/bash {{systemd_notifier_email_folder}}systemd_notifier_email.sh %i
User=root
Group=systemd-journal
