# roles/health-btrfs/templates/health-btrfs.timer.j2
[Unit]
Description=starts health-btrfs.service

[Timer]
OnCalendar={{on_calendar_health_btrfs}}
RandomizedDelaySec={{randomized_delay_sec}}
Persistent=false

[Install]
WantedBy=timers.target
