# roles/certbot-nginx/templates/certbot.timer.j2
[Unit]
Description=Renewal of Let's Encrypt's certificates

[Timer]
OnCalendar=0/12:00:00
RandomizedDelaySec={{randomized_delay_sec}}
Persistent=true

[Install]
WantedBy=timers.target
