# roles/server_native-docker-compose-restart-unhealthy/templates/docker-compose-restart-unhealthy.service.j2
[Unit]
Description=restart unhealthy docker containers
OnFailure=systemd-email@%n.service

[Service]
Type=oneshot
ExecStart=/bin/python {{docker_docker_compose_restart_unhealthy}}server_docker-compose-restart-unhealthy.py
