# roles/native-pull-remote-backups/templates/pull-remote-backups.timer
[Unit]
Description=pulls a remote backup

[Timer]
OnCalendar=22:00
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
