# roles/srv-web-6-6-tls-renew/templates/srv-web-6-6-tls-renew.service.j2
[Unit]
Description=Let's Encrypt renewal
OnFailure=sys-alm-compose.cymais@%n.service

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew --quiet --agree-tos
ExecStartPost=/bin/systemctl reload nginx.service
