# roles/driver-msi-keyboard-color/templates/keyboard-color.service.j2
[Unit]
Description=Keyboard Color Service
OnFailure=systemd-notifier.cymais@%n.service

[Service]
Type=oneshot
ExecStart=/bin/python /opt/keyboard_color.py {{ vendor_and_product_id }}