# roles/driver-msi-keyboard-color/templates/keyboard-color.timer.j2
[Unit]
Description=Keyboard Color Timer

[Timer]
OnCalendar={{on_calendar_msi_keyboard_color}}
Persistent=true

[Install]
WantedBy=timers.target
