# roles/native-disc-space-check/templates/disc-space-check.timer.j2
[Unit]
Description=starts disc-space-check.service

[Timer]
OnCalendar={{on_calendar_disc_space_check}}

[Install]
WantedBy=timers.target
